# app/50-format-usb-storage.rules
ACTION=="add", \
SUBSYSTEMS=="usb", \
RUN+="/opt/usb-storage-formatter/on-usb-plugged-in.sh"
