# app/50-format-usb-memory.rules
ACTION=="add", \
SUBSYSTEMS=="usb", \
RUN+="/bin/bash /opt/usb-memory-formatter/log.sh"
