# app/50-format-usb-drive.rules
ACTION=="add", \
SUBSYSTEMS=="usb", \
DEVTYPE=="disk", \
ID_USB_DRIVER=="usb-storage", \
RUN+="/opt/usb-storage-formatter/on-usb-plugged-in.sh"
